# init/upstart/lifecycled.conf
start on startup
kill timeout 30
respawn
script
    set -a
    . /etc/lifecycled
    exec /usr/bin/lifecycled >> /var/log/lifecycled 2>&1
end script
